# deploy/bigisland.service
[Unit]
Description=Big Island Webapp (Kona Trading Dashboard)
After=network.target

[Service]
User=ec2-user
Group=ec2-user
WorkingDirectory=/home/ec2-user/bigisland
EnvironmentFile=/home/ec2-user/bigisland/.env
ExecStart=/home/ec2-user/.local/bin/gunicorn wsgi:app --bind 127.0.0.1:8001 --workers 1 --threads 4 --timeout 300
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
